# deploy/systemd/secmon-fetch-secrets.service
# Concept: Mukesh Kesharwani
# Contact: mukesh.kesharwani@adobe.com
[Unit]
Description=Refresh Security Monitoring secrets from AWS Secrets Manager
After=network-online.target
Before=secmon-streamlit.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/secmon/app/scripts/ec2_fetch_secrets.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
